# systemd/duckdns.service
[Unit]
Description=DuckDNS IP-change updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duckdns-update.sh
