# systemd/duckdns-force.service
[Unit]
Description=DuckDNS daily refresh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duckdns-force-update.sh
