# systemd/duckdns.timer
[Unit]
Description=DuckDNS IP change checker

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
