# systemd/duckdns-force.timer
[Unit]
Description=DuckDNS daily refresh timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
